# files/eventstore_upstart.conf
description "EventStore"

start on startup
start on runlevel [2345]

stop on shutdown
respawn
exec /vagrant/start_eventstore
